# scripts/price-tracker.service
[Unit]
Description=Price Tracker Telegram Bot
After=network.target price-tracker-xvfb.service price-tracker-fluxbox.service
Requires=price-tracker-xvfb.service price-tracker-fluxbox.service

[Service]
Type=simple
User=yury
WorkingDirectory=/home/yury/price_tracker
Environment=DISPLAY=:99
EnvironmentFile=-/home/yury/price_tracker/.env
ExecStart=/home/yury/price_tracker/venv/bin/python main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
